# .docker/cron/conf/crontab
*/1 * * * * /usr/local/bin/php __working_dir/bin/magento cron:run
*/1 * * * * /usr/local/bin/php __working_dir/update/cron.php
*/1 * * * * /usr/local/bin/php __working_dir/bin/magento setup:cron:run